# project-aws-autoscaling/nodepurge.service
[Unit]
Description=Run mycommand at shutdown
Requires=network.target
DefaultDependencies=no
Before=shutdown.target reboot.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /home/$USER/nodepurge/startupscript.sh
ExecStop=/bin/bash /home/$USER/nodepurge/purgescript.sh

[Install]
WantedBy=multi-user.target
